# tools/vps-setup/systemd/dcc-locale-prewarm.service
[Unit]
Description=HESEM MOM DCC English locale prewarm
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=www-data
Group=www-data
WorkingDirectory=/var/www/eqms.hesemeng.com
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
Environment=DCC_TRANSLATION_WORKER_SLOTS=4
Environment=DCC_TRANSLATION_COMMAND_TIMEOUT_SECONDS=1800
Environment=DCC_TRANSLATION_JOB_MAX_ATTEMPTS=3
TimeoutStartSec=12h
ExecStart=/usr/bin/php8.5 /var/www/eqms.hesemeng.com/tools/scripts/translation/dcc_locale_backfill.php --fpm-env=/etc/php/8.5/fpm/pool.d/mom.conf --only-missing-job --no-spawn-per-job --start-workers=4 --wait-for-workers --worker-slots=4 --max-queue=600 --actor=system.locale-prewarm
